# etc/supervisor.conf

[program:licensedb]
command=/home/site/licensedb/bin/server.sh
user=site
redirect_stderr=true

