# etc/service/log/run
#!/bin/bash

mkdir --parents main
exec setuidgid site multilog t ./main

